# schedule.cron
* 8 * * * python3 /home/opc/OCI-ShowBackups/OCI-ShowBackups.py